# debian/postinst
#!/bin/sh
# postinst script for ctxos-core
set -e

case "$1" in
    configure)
        echo "Configuring CTX OS Core settings..."

        # 1. PostgreSQL: Include the custom OS config
        # We look for any installed PG version and append the include line
        PG_CONF_ROOT="/etc/postgresql"
        if [ -d "$PG_CONF_ROOT" ]; then
            for version in $(ls "$PG_CONF_ROOT"); do
                CONF_FILE="$PG_CONF_ROOT/$version/main/postgresql.conf"
                if [ -f "$CONF_FILE" ]; then
                    if ! grep -q "ctx_postgresql.conf" "$CONF_FILE"; then
                        echo "include = '/etc/postgresql/ctx_postgresql.conf'" >> "$CONF_FILE"
                        echo "PostgreSQL $version tuned."
                    fi
                fi
            done
        fi

        # 2. Systemd: Enable and Start ctxos-dns
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload >/dev/null || true
            if [ -f "/lib/systemd/system/ctxos-dns.service" ]; then
                systemctl enable ctxos-dns.service >/dev/null || true
                systemctl restart ctxos-dns.service >/dev/null || true
            fi
        fi

        # 2.1. Plymouth: Set boot splash theme
        if [ -x "/usr/sbin/plymouth-set-default-theme" ]; then
            if [ -d "/usr/share/plymouth/themes/ctx-pulse" ]; then
                plymouth-set-default-theme ctx-pulse >/dev/null 2>&1 || true
                update-initramfs -u >/dev/null 2>&1 || true
            fi
        fi

        # 3. Hardware: Reload Udev rules for HyperFIDO
        if [ -x "/usr/bin/udevadm" ]; then
            udevadm control --reload-rules || true
            udevadm trigger || true
        fi

        # 4. Networking: Ensure Samba public directory exists
        if [ -f "/etc/samba/smb.conf" ]; then
            mkdir -p /home/public
            chmod 777 /home/public
        fi

        # 5. Apt: Set permissions for sensitive files
        chmod 644 /etc/apt/preferences.d/ctx-pinning || true
        chmod 644 /etc/apt/sources.list.d/ctx.list || true

        # 6. GPG key for CTX OS repository
        # Note: GPG key should be included in the package (debian/install-root/usr/share/keyrings/)
        # The key file should be placed at /usr/share/keyrings/ctxos-archive-key.gpg
        if [ ! -f /usr/share/keyrings/ctxos-archive-key.gpg ] && [ ! -f /etc/apt/keyrings/ctxos-archive-key.gpg ]; then
            echo "WARNING: CTX OS GPG key not found. Please install it manually."
        fi

        # 7. Desktop: Update icon cache if you added any
        if [ -x "/usr/bin/update-desktop-database" ]; then
            update-desktop-database -q
        fi

        # 8. Update dconf database for system-wide GNOME settings
        if [ -x "/usr/bin/dconf" ]; then
            dconf update || true
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.
#DEBHELPER#

exit 0
